# misc/hypercorn.conf
[program:hypercorn]
command=/usr/local/bin/python main.py
directory=/deploy/wiki_app